# systemd/ddns.service
[Unit]
Description=Update Cloudflare dynamic dns

[Service]
Type=simple
Restart=no
ExecStart=/opt/ddns.sh
